# torrentd/nimos-torrentd.service
[Unit]
Description=NimOS Torrent Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nimos-torrentd
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5
User=nimbus
Group=nimbus
RuntimeDirectory=nimos
StateDirectory=nimos/torrentd

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data /var/lib/nimos /run/nimos /nimbus/pools

[Install]
WantedBy=multi-user.target
